# crontab
* * * * * root php /var/apps/grd-crawler/current/jobby.php >> /var/apps/grd-crawler/current/logs/jobby.log 2>&1